# service_files/xyz.openbmc_project.nvmesensor.service
[Unit]
Description=NVMe Sensor
StopWhenUnneeded=false
After=xyz.openbmc_project.FruDevice.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/sbin/nvmesensor

[Install]
WantedBy=multi-user.target